# build/deb/DEBIAN/postinst
#!/bin/sh
set -e

RIC_HOME=/var/lib/ric
RIC_USER=ric

adduser --system --group --home "$RIC_HOME" "$RIC_USER"

systemctl enable ric
systemctl start ric
