# build/deb/DEBIAN/prerm
#!/bin/sh
set -e

RIC_HOME=/var/lib/ric
RIC_USER=ric

systemctl stop ric
systemctl disable ric

rm /usr/local/bin/ric
